# oscar-3.1.11/usr/GNUstep/.config/autostart/Terminal.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=OSCAR
Comment=OSCAR
Exec=/usr/share/oscar/bin/lancer_oscar_dans_terminal
Hidden=false